# skilling-continuation-labs/living-off-the-land-attacks/lab/skillshub/logInject.timer
[Unit]
Description=Timer for logInject Service

[Timer]
OnBootSec=5min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target