# pc6-round1/team/exfil-enigma/challenge/wan-1/netcat-listener.service
[Unit]
Description=Netcat Listener for PDF Transfer
After=network.target

[Service]
ExecStart=/home/user/Documents/netcatlistener.sh
Restart=always

[Install]
WantedBy=multi-user.target