# pc4/individual-b/round2-why-is-the-rd-always-gone/challenge/ssh-server/add-ssh-user.service
[Unit]
Description=Execute adding a user

[Service]
ExecStart=/etc/systemd/system/add-ssh-user.sh
Restart=always
RestartSec=300

[Install]
WantedBy=default.target
