# pc4/individual-b/round2-why-is-the-rd-always-gone/challenge/ssh-server/api-placement.service
[Unit]
Description=Execute adding a file with API key

[Service]
ExecStart=/etc/systemd/system/api-placement.sh

[Install]
WantedBy=default.target
