# pc4/individual-b/round2-why-is-the-rd-always-gone/challenge/zulip-server/startup-api.service
[Unit]
Description=Execute system to add API key
ConditionPathExists=!/etc/systemd/system/startup-api-done.txt

[Service]
Type=simple
ExecStart=/etc/systemd/system/startup-api.sh

[Install]
WantedBy=default.target
