# pc5/individual-b/round1-nice-encryption-you-have-there/challenge/build/code-server.service
# Copyright 2024 Carnegie Mellon University.
# Released under a MIT (SEI)-style license, please see LICENSE.md in the project 
# root or contact permission@sei.cmu.edu for full terms.

[Unit]
Description=Code Server Startup
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/code
ExecStart=/home/user/code/target/release/server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
