# pc4/individual-a/round3-oops-we-did-it-again/challenge/compromised/p3rs1st.timer
[Unit]
Description=checker

[Timer]
OnCalendar=*-*-* *:*:00/5
Persistent=true

[Install]
WantedBy=timers.target


