# pc7/individual-offense/round1-code_osiris/challenge/code_osiris_remote/code_osiris_remote.socket
[Unit]
Description=Code Osiris Remote Test TCP Socket

[Socket]
ListenStream=9999

[Install]
WantedBy=sockets.target