# pc7/individual-offense/round1-code_osiris/challenge/code_osiris_abyssnet/abyssnet.socket
[Unit]
Description=Code Osiris Abyssnet Socket

[Socket]
ListenStream=9999

[Install]
WantedBy=sockets.target
