# pc4/individual-b/round2-why-is-the-rd-always-gone/challenge/ubuntu-user-sim/router-add-user.service
[Unit]
Description=Execute adding user to router on startup

[Service]
ExecStart=/etc/systemd/system/router-add-user.sh
Restart=always
RestartSec=60

[Install]
WantedBy=default.target
