# pc4/individual-b/round2-why-is-the-rd-always-gone/challenge/ubuntu-user-sim/zulip-poster.service
[Unit]
Description=Execute adding chat to Zulip on startup

[Service]
ExecStart=/etc/systemd/system/zulip-poster.sh

[Install]
WantedBy=default.target
